# etc/cron.d/com-cssminifier
MAILTO=andychilton@gmail.com

# mins hour dofm mnth dofw   user          command
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
  */10    *    *    *    *   chilts        find /var/lib/com-cssminifier/ -mmin +10 -delete
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
